# nixos-modules/image/installer/inittab
::sysinit:/etc/init.d/rcS
tty2::askfirst:-/bin/sh
tty3::askfirst:-/bin/sh
tty4::askfirst:-/bin/sh
tty5::askfirst:-/bin/sh
::restart:/sbin/init
null::ctrlaltdel:/sbin/reboot
null::shutdown:/bin/umount -a -r
null::shutdown:/sbin/swapoff -a

# NOTE: Busybox will run these in order, so ensure that systemd-udevd is
# started before the installer.
null::respawn:/opt/systemd/lib/systemd/systemd-udevd
::respawn:/opt/installer/bin/installer
